# playbooks/files/rc.local
# disable touchscreen display
xinput disable `xinput|grep ELAN24EE|sed 's/.*id\=//'|sed 's/\t.*$//'`

